# docker/cron/crontab
*/5 * * * * php /var/www/parsing_news_site/bin/console parse_news
# Don't remove the empty line at the end of this file. It is required to run the cron job
